# apps/group-project-gg-bond-main/backend/ggbond-mysql-backup.service
[Unit]
Description=GG Bond MySQL Backup
After=mysql.service
Requires=mysql.service

[Service]
Type=oneshot
User=destiny
WorkingDirectory=/home/destiny/backend
ExecStart=/home/destiny/backend/backup-mysql.sh
